# root/etc/s6-overlay/s6-rc.d/svc-mod-transmission-seed-time/run
#!usr/bin/with-contenv bash

python3 /etc/s6-overlay/s6-rc.d/svc-mod-transmission-seed-time/run.py