# root/etc/s6-overlay/s6-rc.d/init-mod-transmission-seed-time-postinstall/run
#!usr/bin/with-contenv bash

echo "hello"
python3 /etc/s6-overlay/s6-rc.d/init-mod-transmission-seed-time-postinstall/run.py